# delete_file.list
/usr/lib/systemd/user/lipstick.service.d/50-compositor.conf
/usr/lib/systemd/user/lipstick.service.d/50-vendor.hwcomposer-2-3.conf
/usr/lib/systemd/system/droid-bootctl.service
/usr/bin/droid/droid-bootctl.sh
